# dokku/monitoring/memory-exporter/memory_exporter.service
[Unit]
Description=RAM memory exporter to prometheus pushgateway

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/memory_exporter.sh